# shared/systemd/penpot-compose.service
[Unit]
Description=Penpot Docker Compose Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/penpot/compose
EnvironmentFile=/etc/penpot/penpot.env
ExecStart=/usr/bin/flock -w 900 /run/penpot-compose.lock /usr/bin/docker compose -f /opt/penpot/compose/docker-compose.yaml up -d
ExecStop=/usr/bin/flock -w 900 /run/penpot-compose.lock /usr/bin/docker compose -f /opt/penpot/compose/docker-compose.yaml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
